# supervisor-config.conf
# Supervisor configuration for Lisa Sales Agent
# File: /etc/supervisor/conf.d/lisa-sales-agent.conf

[program:lisa-sales-agent]
command=/home/lisaagent/lisa-sales-agent/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
directory=/home/lisaagent/lisa-sales-agent/backend
user=lisaagent
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/lisa-sales-agent.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=
    OPENAI_API_KEY="your_openai_api_key_here",
    SECRET_KEY="your_super_secret_jwt_key",
    GOOGLE_CRED_JSON="/home/lisaagent/lisa-sales-agent/Google_Credentials.json",
    SHEET_NAME="Your Google Sheet Name",
    ENVIRONMENT="production",
    CONVERSATION_SCRIPT_PATH="/home/lisaagent/lisa-sales-agent/conversation_script.txt"
